# systemd/razer-naga-2014-key-remap.service
[Unit]
Description=remap for razer naga 2014 keys

[Service]
Type=simple
ExecStart=/usr/bin/razer-naga-2014-key-remap

[Install]
WantedBy=multi-user.target
